# scripts/utils/simulator_bank/tor_simulador_onion.service
[Unit]
Description=Tor Hidden Service - Simulador Bancario
After=network.target tor.service
Requires=tor.service

[Service]
Type=oneshot
ExecStart=/root/configurar_tor_onion.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
